# deploy/vm/crontab
# Stock Research — production schedule (all times IST / Asia/Kolkata)
SHELL=/bin/bash
TZ=Asia/Kolkata
PATH=/opt/stock-research/venv/bin:/usr/bin:/bin
PYTHONPATH=/opt/stock-research

# Daily research run — 06:30 IST Mon-Fri
30 6 * * 1-5  cd /opt/stock-research && python run_agents.py research >> /var/log/stock-research.log 2>&1

# Evening intraday watchlist — 18:30 IST Mon-Fri
30 18 * * 1-5  cd /opt/stock-research && python run_agents.py intraday >> /var/log/stock-intraday.log 2>&1

# Merged watcher every 3 min Mon-Fri (self-skips outside 09:15-15:30 IST)
*/3 * * * 1-5  cd /opt/stock-research && python run_agents.py watch >> /var/log/stock-watch.log 2>&1
